# marvin-ts/src/wakeword/marvin-wakeword.service
[Unit]
Description=Marvin Wake Word Detector (Hey Marvin)
Documentation=https://github.com/dscripka/openWakeWord
After=sound.target pulseaudio.service pipewire.service
Wants=pulseaudio.service

[Service]
Type=simple
ExecStart=/usr/bin/sg audio -c "%h/copilot-assistant-thing/marvin-ts/.venv/bin/python %h/copilot-assistant-thing/marvin-ts/src/wakeword/wakeword.py"
WorkingDirectory=%h/copilot-assistant-thing/marvin-ts
Restart=on-failure
RestartSec=5

# Node.js must be on PATH for headless Marvin queries (exclude conda to avoid lib conflicts)
Environment=PATH=%h/.nvm/versions/node/v22.20.0/bin:/usr/local/bin:/usr/bin:/bin
# Work around conda libstdc++ / libgcc_s conflicts
Environment=LD_PRELOAD=/usr/lib/x86_64-linux-gnu/libstdc++.so.6
# Force ALSA to use the ReSpeaker as default capture device
Environment=ALSA_CARD=ArrayUAC10

# Tuning (uncomment and adjust as needed)
# Environment=MARVIN_WAKE_THRESHOLD=0.5
# Environment=MARVIN_WAKE_COOLDOWN=3.0
# Environment=MARVIN_WAKE_MODEL=/path/to/custom_model.onnx

# Inherit user session for audio access
Environment=DISPLAY=:0
PassEnvironment=DISPLAY PULSE_SERVER XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS

[Install]
WantedBy=default.target
